# restore-agents.timer
[Unit]
Description=Timer for restoring agents after boot
After=agent-dashboard.service

[Timer]
OnBootSec=30s
OnUnitActiveSec=1h
Persistent=true

[Install]
WantedBy=timers.target